# .config/supervisor/uwsgi.conf
[program:uwsgi]
startsecs = 0
autorestart = false
command=/usr/local/bin/uwsgi --ini /tmp/.config/uwsgi.ini